# deploy/trading-deploy.timer
[Unit]
Description=Run trading auto-deploy every 2 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
